# systemd/hdmi-edid-user.service
[Unit]
Description=HDMI EDID User Session Refresh
After=default.target
Wants=default.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'systemctl --user restart pipewire pipewire-pulse wireplumber; sleep 1; killall wf-panel-pi'
RemainAfterExit=yes

[Install]
WantedBy=default.target
